# roles/b2-sync/templates/b2-sync.timer
[Unit]
Description=B2 Sync

[Timer]
OnCalendar=Sat *-*-* 05:00:00

[Install]
WantedBy=timers.target
